# examples/goflow2-aggregator.service
[Unit]
Description=goflow2 NetFlow/IPFIX Aggregator
After=network.target
Wants=network.target

[Service]
Type=simple
User=goflow
Group=goflow
WorkingDirectory=/opt/goflow2-aggregator

ExecStart=/bin/bash -c '/usr/local/bin/goflow2 -listen netflow://:2055 | /opt/goflow2-aggregator/goflow2-aggregator'

Restart=always
RestartSec=10

Environment="RUST_LOG=info"

StandardOutput=journal
StandardError=journal
SyslogIdentifier=goflow2-aggregator

[Install]
WantedBy=multi-user.target
